# docs/hook.conf
description "Hook Script"
author "Julien C."

start on runlevel [345]
stop on runlevel [!345]

respawn

exec python /var/hook.py
